# systemd/kstor.service
[Unit]
Description=KStor server

[Service]
Type=notify
ExecStart=kstor-srv --config /path/to/kstor/server-config.yaml
KillSignal=SIGINT
