# systemd/kstor.socket
[Unit]
Description=KStor listen socket

[Socket]
ListenStream=/run/kstor.socket

[Install]
WantedBy=sockets.target
